# etc/cron.d/clean-hls-dir
*/10 * * * * /clean-hls-dir.sh  > /dev/stdout